# holidaylights_stop.service
[Unit]
Description=Stop Holiday Lights Service

[Service]
Type=oneshot
ExecStart=/bin/systemctl stop holiday-lights.service