# holidaylights_stop.timer
[Unit]
Description=Stop holiday lights at 7 AM

[Timer]
OnCalendar=*-*-* 07:00:00
Persistent=true
Unit=holiday-lights-stop.service

[Install]
WantedBy=timers.target